# os/OSMC/custom/etc/udev/rules.d/usbsound.rules
ACTION=="add", KERNEL=="controlC1", SUBSYSTEM=="sound", RUN+="/home/osmc/soundusb.py"
ACTION=="remove", KERNEL=="controlC1", SUBSYSTEM=="sound", RUN+="/home/osmc/soundhdmi.py"
